# tele_clone_saham.service
[Unit]
Description=Tele Clone Saham forwarder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zero
WorkingDirectory=/home/zero/Documents/BOT/tele_clone_saham
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/home/zero/Documents/BOT/tele_clone_saham/.env
ExecStart=/home/zero/Documents/BOT/tele_clone_saham/venv/bin/python /home/zero/Documents/BOT/tele_clone_saham/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
